# server/sipwitch.cron
#!/bin/sh
if test -p /var/run/sipwitch/control ; then
        sipwitch period 60 ; fi

